# infra/launchd/com.hoya.jangmanal.crawler.emart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hoya.jangmanal.crawler.emart</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/hoya/workspace/JangMaNal/apps/crawler/scripts/run_crawler_daily.sh</string>
    <string>emart</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/hoya/workspace/JangMaNal/apps/crawler</string>
  <key>StandardOutPath</key>
  <string>/Users/hoya/workspace/JangMaNal/apps/crawler/logs/launchd-emart.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hoya/workspace/JangMaNal/apps/crawler/logs/launchd-emart.err.log</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
